# deploy/autodeploy/noblechat-autodeploy.timer
[Unit]
Description=Poll GitHub main and redeploy NobleChat

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=15s

[Install]
WantedBy=timers.target
