# deploy/autodeploy/noblechat-autodeploy.service
[Unit]
Description=NobleChat auto-deploy (redeploy when main advances)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/noblechat-autodeploy.sh
